# .config/systemd/user/update-git@.timer
[Unit]
Description=Update git repository %I

[Timer]
OnCalendar=Daily
Persistent=true
Unit=update-git@%i.service

[Install]
WantedBy=timers.target
